# systemd/.config/systemd/user/kanshi.service
[Unit]
Description=Kanshi service
Documentation=man:kanshi(1)

[Service]
Type=simple
ExecStart=/usr/bin/kanshi
